# docker/supervisor.conf
[program:redis-server]
command=redis-server
autorestart=true

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true

[program:flask]
command=./app.py
directors=/app/src/backend
